# dns-updater.service
[Unit]
Description=Cloudflare DNS Updater Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
EnvironmentFile=/etc/cfdns-updater.env
WorkingDirectory=/opt/cfdns-updater
ExecStart=/opt/cfdns-updater/cfdns-updater
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
